# launchagent/com.luciffer.claude-dual-rotate-audit.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.luciffer.claude-dual-rotate-audit</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/luciffer/.claude-dual/rotate-audit-log.sh</string>
    </array>
    <!-- Run once a day at 03:17 local time (off the :00/:30 cliff). -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>17</integer>
    </dict>
    <!-- Also run once on load (catches missed runs after sleep/reboot). -->
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/luciffer/.claude-dual/rotate-audit.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/luciffer/.claude-dual/rotate-audit.err.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
